# tcm-diagnosis.service
[Unit]
Description=TCM Intelligent Diagnosis System
After=network.target

[Service]
Type=simple
User=www-data
WorkingDirectory=/www/wwwroot/aiagents-zhongyi/backend
Environment="PATH=/www/wwwroot/aiagents-zhongyi/venv/bin"
ExecStart=/www/wwwroot/aiagents-zhongyi/venv/bin/uvicorn main:app --host 0.0.0.0 --port 8000
Restart=always
RestartSec=10

[Install]
WantedBy=multi-user.target

